# cookbooks/blockstep/templates/default/foreman/master.conf.erb
pre-start script

bash << "EOF"
EOF

end script


start on (net-device-up
          and local-filesystems
          and runlevel [2345])

stop on [!2345]

post-stop script

sleep 2

end script
